# bobstudio-simple.service
[Unit]
Description=BOB Studio Node.js Application
Documentation=https://github.com/sunshaoxuan/bobstudio
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bobstudio
Environment="NODE_ENV=production"
# 加载用户的环境变量
Environment="PATH=/root/.nvm/versions/node/v18.20.5/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# 如果使用 nvm，需要加载 nvm
ExecStartPre=/bin/bash -lc 'cd /root/bobstudio && git pull origin main'
ExecStartPre=/bin/bash -lc 'cd /root/bobstudio && npm install'
ExecStartPre=/bin/bash -lc 'cd /root/bobstudio && npm run build'
ExecStart=/bin/bash -lc 'cd /root/bobstudio && npm run server'
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/bobstudio/output.log
StandardError=append:/var/log/bobstudio/error.log

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

