# bobstudio-simple.service
[Unit]
Description=BOB Studio Node.js Application
Documentation=https://github.com/sunshaoxuan/bobstudio
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bobstudio
Environment="NODE_ENV=production"
ExecStart=/bin/bash /root/bobstudio/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/bobstudio/output.log
StandardError=append:/var/log/bobstudio/error.log

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
