# bobstudio.service
[Unit]
Description=BOB Studio Node.js Application
Documentation=https://github.com/sunshaoxuan/bobstudio
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bobstudio
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/bin/bash -c 'cd /root/bobstudio && /usr/bin/git pull origin main'
ExecStartPre=/bin/bash -c 'cd /root/bobstudio && /usr/bin/npm install'
ExecStartPre=/bin/bash -c 'cd /root/bobstudio && /usr/bin/npm run build'
ExecStart=/bin/bash -c 'cd /root/bobstudio && /usr/bin/npm run server'
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/bobstudio/output.log
StandardError=append:/var/log/bobstudio/error.log

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

